# ics-sync.timer
[Unit]
Description=Poll for new calendar emails every 2 minutes
After=proton-bridge.service

[Timer]
OnBootSec=60
OnUnitActiveSec=120
AccuracySec=15

[Install]
WantedBy=timers.target
